# files/mukube-configurator/etc/systemd/system/mnt.mount
[Unit]
StopWhenUnneeded=true

[Mount]
Where=/mnt
What=/dev/disk/by-label/config
Type=ext4
Options=defaults
