# files/misc/etc/systemd/system/systemd-journal-flush.service.d/override.conf
# https://github.com/systemd/systemd/commit/fe7f113c1331e11af4e9f815d7305c8b1b332312
[Unit]
ConditionPathExists=!/etc/initrd-release
